# layers/meta-quectel/recipes-quectel/tests/files/rgpiod.service
[Unit]
Description=rgpiod deamon

[Service]
Type=simple
ExecStart=/usr/bin/rgpiod

[Install]
WantedBy=multi-user.target
