# layers/meta-quectel/recipes-quectel/quectel/files/quectel-cm@.service
[Unit]
Description=quectel connection manager tool

[Service]
Type=simple
ExecStart=/usr/bin/quectel-CM -i %I

[Install]
